# deploy/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:backend]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8001 --workers 4
directory=/app/backend
user=quantum
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/backend.err.log
stdout_logfile=/var/log/supervisor/backend.out.log
environment=PYTHONPATH="/app:/app/backend:/app/crypto-core:/app/ai-engine:/app/blockchain-integration"

[program:nginx]
command=nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
priority=999

[program:mongodb]
command=mongod --dbpath /data/db --port 27017
user=mongodb
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/mongodb.err.log
stdout_logfile=/var/log/supervisor/mongodb.out.log
priority=100

[program:ai-monitor]
command=python -m ai_engine.monitoring.real_time_monitor
directory=/app
user=quantum
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/ai-monitor.err.log
stdout_logfile=/var/log/supervisor/ai-monitor.out.log
environment=PYTHONPATH="/app:/app/ai-engine"

[program:threat-scanner]
command=python -m ai_engine.threat_detection.threat_scanner
directory=/app
user=quantum
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/threat-scanner.err.log
stdout_logfile=/var/log/supervisor/threat-scanner.out.log
environment=PYTHONPATH="/app:/app/ai-engine"

[group:quantumgate]
programs=backend,nginx,mongodb,ai-monitor,threat-scanner
priority=999

[inet_http_server]
port=127.0.0.1:9001

[supervisorctl]
serverurl=http://127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface